# deploy/stephenasatsa-go-password.service
[Unit]
Description=Stephen Asatsa Password Service (Go)
After=network.target

[Service]
Type=simple
User=http
Group=http
WorkingDirectory=/opt/stephenasatsa/apps/website/backend/go-services
Environment=PORT=9001
ExecStart=/opt/stephenasatsa/apps/website/backend/go-services/password-service/password-service
Restart=always
RestartSec=10
StandardOutput=append:/var/log/stephenasatsa/go-password.log
StandardError=append:/var/log/stephenasatsa/go-password-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stephenasatsa

[Install]
WantedBy=multi-user.target
