# deploy/stephenasatsa-go-image.service
[Unit]
Description=Stephen Asatsa Image Service (Go)
After=network.target

[Service]
Type=simple
User=http
Group=http
WorkingDirectory=/opt/stephenasatsa/apps/website/backend/go-services
Environment=PORT=9003
ExecStart=/opt/stephenasatsa/apps/website/backend/go-services/image-service/image-service
Restart=always
RestartSec=10
StandardOutput=append:/var/log/stephenasatsa/go-image.log
StandardError=append:/var/log/stephenasatsa/go-image-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stephenasatsa

[Install]
WantedBy=multi-user.target
